# sv/tt-rss-update/run
#!/bin/sh
cd /opt/Tiny-Tiny-RSS-1.15.3
exec chpst -u www-data /opt/Tiny-Tiny-RSS-1.15.3/update_daemon2.php
