# sv/tt-rss-update/run
#!/bin/sh
cd /opt/Tiny-Tiny-RSS-1.13
exec setuidgid www-data /opt/Tiny-Tiny-RSS-1.13/update_daemon2.php
