# sv/tt-rss-update/run
#!/bin/sh
cd /opt/Tiny-Tiny-RSS-master
exec chpst -u www-data /opt/Tiny-Tiny-RSS-master/update_daemon2.php
