# sv/tt-rss-update/run
#!/bin/sh
cd /opt/tt-rss.git
exec chpst -u www-data /opt/tt-rss.git/update_daemon2.php
